# host-image/lib/systemd/system/reflect-passwd@.service
[Unit]
Description=Reflect Host Password and Group Info into %i Container

[Service]
ExecStart=/usr/sbin/reflect-passwd %i
