# host-image/lib/systemd/system/reflect-apps@.path
[Unit] 
Description=Watch for App Changes in the %i Container

[Path]
PathExists=/var/lib/machines/%i/usr/share/applications
PathModified=/var/lib/machines/%i/usr/share/applications

[Install]
WantedBy=multi-user.target

