# host-image/lib/systemd/system/reflect-locale@.path
[Unit] 
Description=Watch for Changes to Host Locale

[Path]
PathExists=/etc/default/locale
PathChanged=/etc/default/locale

[Install]
WantedBy=multi-user.target
