# host-image/lib/systemd/system/reflect-apps@.service
[Unit]
Description=Reflect %i Container's Apps as Launchers on Host

[Service]
ExecStart=/usr/sbin/reflect-apps %i
