# host-image/lib/systemd/system/reflect-timezone@.path
[Unit] 
Description=Watch for Changes to Host Timezone

[Path]
PathExists=/etc/timezone
PathChanged=/etc/timezone

[Install]
WantedBy=multi-user.target

