# host-image/lib/systemd/system/reflect-timezone@.service
[Unit]
Description=Reflect Host Timezone into %i Container

[Service]
ExecStart=/usr/sbin/reflect-timezone %i
