# host-image/usr/share/ds64/container-image/lib/systemd/system/reflect-timezone.service
[Unit]
Description=Reflect Host Timezone into Container

[Service]
ExecStart=/usr/sbin/reflect-timezone

