# host-image/usr/share/ds64/container-image/lib/systemd/system/reflect-locale.path
[Unit] 
Description=Watch for Changes to Host Locale

[Path]
PathExists=/etc/default/host-locale
PathChanged=/etc/default/host-locale

[Install]
WantedBy=multi-user.target
