# host-image/usr/share/ds64/container-image/lib/systemd/system/reflect-timezone.path
[Unit] 
Description=Watch for Changes to Host Timezone

[Path]
PathExists=/etc/host-timezone
PathChanged=/etc/host-timezone

[Install]
WantedBy=multi-user.target

