# host-image/usr/share/ds64/container-image/lib/systemd/system/reflect-locale.service
[Unit]
Description=Reflect Host Locale into Container

[Service]
ExecStart=/usr/sbin/reflect-locale
