# etc-systemd-system/reflect-passwd.path
[Unit] 
Description=reflect-passwd

[Path]
PathExists=/etc/passwd
PathModified=/etc/passwd
PathModified=/etc/group
PathModified=/etc/shadow
PathModified=/etc/gshadow

[Install]
WantedBy=multi-user.target

