# etc-systemd-system/reflect-apps.service
[Unit]
Description=reflect-apps

[Service]
ExecStart=/usr/local/bin/reflect-apps

