# etc-systemd-system/reflect-apps.path
[Unit] 
Description=reflect-apps

[Path]
PathExists=/var/lib/machines/debian-stretch-64/usr/share/applications
PathModified=/var/lib/machines/debian-stretch-64/usr/share/applications

[Install]
WantedBy=multi-user.target

